# LoRa/kubernetes/kubectl-portforward-lora.service
[Unit]
Description=Kubectl port-forward for lora-receiver
After=network.target k3s.service
Requires=k3s.service

[Service]
Type=simple
User=pi
Environment=KUBECONFIG=/etc/rancher/k3s/k3s.yaml
ExecStart=/usr/local/bin/kubectl port-forward -n lora-demo svc/lora-receiver 8080:8080
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
